# packer/linux/stack/conf/buildkite-agent/systemd/buildkite-agent-health-check-onfailure.service
[Unit]
Description=Health check for buildkite-agent triggered by service failure

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-agent-health --skip-uptime-check
