# packer/linux/conf/buildkite-agent/systemd/cloud-final.service.d/10-power-off-on-failure.conf
[Unit]
OnFailure=poweroff.target
OnFailureJobMode=replace-irreversibly
