# packer/conf/docker/cron.daily/docker-gc
#!/bin/bash
set -euo pipefail

# prune stopped containers and images
docker system prune -f
