# packer/conf/docker/cron.daily/docker-gc
#!/bin/bash
set -euo pipefail

# prune stopped containers and images older than 7 days
docker system prune –force –filter "until=10080m"