# packer/linux/conf/buildkite-agent/systemd/buildkite-agent-hard-limit.service
[Unit]
Description=Hard stop buildkite agent

[Service]
ExecStart=/bin/systemctl stop buildkite-agent
RemainAfterExit=true
