# packer/linux/conf/buildkite-agent/systemd/buildkite-agent-soft-limit.service
[Unit]
Description=Soft stop buildkite agent

[Service]
ExecStart=/bin/systemctl kill --signal=TERM --kill-who=main buildkite-agent
RemainAfterExit=true
