# packer/linux/conf/buildkite-agent/systemd/buildkite-agent-soft-limit.timer
[Unit]
Description=Initiate soft stop of agent
After=buildkite-agent.service

[Timer]
OnBootSec=3hr

[Install]
WantedBy=timers.target
