# packer/conf/buildkite-lifecycle-agent.conf
start on startup
respawn
script
    mkfifo /tmp/buildkite-lifecycle-agent-log-fifo
    ( logger -t buildkite-lifecycle-agent </tmp/buildkite-lifecycle-agent-log-fifo & )
    exec >/tmp/buildkite-lifecycle-agent-log-fifo
    exec 2>&1
    rm /tmp/buildkite-lifecycle-agent-log-fifo
    exec /usr/bin/buildkite-lifecycle-agent -pidfile
end script